# growup-sensor.service
[Unit]
Description=GrowUp IoT Sensor Agent
Documentation=file:///home/pi/growup-python/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/growup-python
Environment="PATH=/home/pi/growup-python/venv/bin"
EnvironmentFile=/home/pi/growup-python/.env

# Start the sensor agent
ExecStart=/home/pi/growup-python/venv/bin/python3 main.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=growup-sensor

# Resource limits
MemoryMax=512M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
