# systemd/growup-iot.service
[Unit]
Description=GrowUp IoT System - Modular Architecture
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/growup-iot
Environment="PATH=/home/pi/growup-iot/venv/bin"
ExecStart=/home/pi/growup-iot/venv/bin/python3 /home/pi/growup-iot/main.py

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=growup-iot

# Security
PrivateTmp=yes
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/growup-iot
SupplementaryGroups=gpio i2c spi video

[Install]
WantedBy=multi-user.target
